# build_tools/github_actions/runner/config/systemd/system/runner-setup.target
[Unit]
Description=Coordinates startup of all services related to the GitHub runner.

[Install]
WantedBy=multi-user.target
